# deploy/systemd/quark-supervisor.service
[Unit]
Description=Quark Supervisor
Documentation=https://github.com/quarkloop/quark
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=quark
Group=quark
EnvironmentFile=/etc/quark/supervisor.env
ExecStart=/usr/local/bin/quark-supervisor start --port ${QUARK_SUPERVISOR_PORT} --nats-state-dir ${QUARK_NATS_STATE_DIR} --nats-client-port ${QUARK_NATS_CLIENT_PORT} --nats-websocket-port ${QUARK_NATS_WEBSOCKET_PORT} --nats-monitor-port ${QUARK_NATS_MONITOR_PORT}
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/quark /var/log/quark

[Install]
WantedBy=multi-user.target
